# daemontools-tutorial/play-scala-setup/root/etc/service/service-example/run
#!/bin/sh
exec envdir env-vars /bin/sh -c '"$APPSDIR/$APPNAME/bin/$APPNAME" -Dconfig.file=/etc/app.conf -Dlogger.file=/etc/app_logback.xml'

