# daemontools-tutorial/spring-boot-setup/root/etc/service/service-example/run
#!/bin/sh
exec envdir env-vars /bin/sh -c 'java -jar "$APPDIR/$APPNAME" --spring.profiles.active=$ENVPROFILE'

